# admin_site/securityscripts/91-detectusbevents.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/etc/bibos/security/securityEvent.py 'Product:' 'Manufacturer:' 'input:'"