# admin_site/securityscripts/91-detectusbevents.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/etc/bibos/security/keyboardEvent.py 'input:'"